# archFiles/airootfs/etc/xdg/autostart/dark-mode.desktop
[Desktop Entry]
Type=Application
Name=ApplyDarkMode
Exec=gsettings set org.gnome.desktop.interface color-scheme 'prefer-dark'
