# archFiles/airootfs/etc/xdg/autostart/gsettings-tweaks-init.desktop
[Desktop Entry]
Type=Application
Name=GnomeSettingsTweaks

Exec=sh /etc/xdg/autostart/gsettings-tweaks.sh