# archFiles/airootfs/etc/xdg/autostart/gsettings-tweaks-init.desktop
[Desktop Entry]
Type=Application
Name=GnomeSettingsTweaks

Exec=sh gsettings-tweaks.sh