# archFiles/airootfs/etc/xdg/autostart/installer.desktop
[Desktop Entry]
Type=Application
Name=Installer
Exec=kgx -e 'sudo archinstall'
